# tacacs+-FB4.0.4.19.1/tac_plus6.sysvinit
#!/bin/bash
#
# /etc/rc.d/init.d/tac_plus6
#
# chkconfig: 2345 86 14
# description: TACACS+ v6 Daemon

# Define variables
TACPLUS_BIN=/usr/bin/tac_plus6
TACPLUS_OPTS="-B ::" 
TACPLUS_CONF=/etc/tac_plus.conf

# Source function library.
. /etc/rc.d/init.d/functions

RETVAL=0
prog="tac_plus6"

case "$1" in
  start)
        # Got Config File?
        [ -f /etc/tac_plus.conf ] || exit 69

        echo -n $"Starting $prog: "
        daemon $TACPLUS_BIN -C $TACPLUS_CONF $TACPLUS_OPTS
        RETVAL=$?
        [ $RETVAL -eq 0 ] && touch /var/lock/subsys/$prog
        echo
        ;;
  stop)
        echo -n $"Shutting down $prog: "
        killproc $prog
        RETVAL=$?
        [ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/$prog
        echo
        ;;
  restart|reload)
        $0 stop
        $0 start
        RETVAL=$?
        ;;
  status)
        status $prog
        RETVAL=$?
        ;;
  *)
        echo $"Usage: $0 {start|stop|restart|reload|status}"
        exit 2
esac

exit $RETVAL
